# profile.d/eclipse.csh
#!/bin/csh
setenv ECLIPSE_HOME /usr/local/share/eclipse
setenv PATH ${PATH}:${ECLIPSE_HOME}
