# profile.d/eclipse.sh
#!/bin/sh
export ECLIPSE_HOME=/usr/local/share/eclipse
export PATH="${PATH}:${ECLIPSE_HOME}"
