# profile.d/squirrel-sql.csh
#!/bin/csh
setenv SQUIRREL_SQL_HOME /usr/local/share/squirrel-sql
setenv PATH ${PATH}:${SQUIRREL_SQL_HOME}
