# profile.d/squirrel-sql.sh
#!/bin/sh
export SQUIRREL_SQL_HOME=/usr/local/share/squirrel-sql
export PATH="${PATH}:${SQUIRREL_SQL_HOME}"
