# profile.d/terraform.sh
#!/bin/sh
export TERRAFORM_HOME=/usr/local/share/terraform
export PATH="${PATH}:${TERRAFORM_HOME}"
