# profile.d/terraform.csh
#!/bin/csh
setenv TERRAFORM_HOME /usr/local/share/terraform
setenv PATH ${PATH}:${TERRAFORM_HOME}
